# football-squares-main.service
[Unit]
Description=Football Squares Main Branch Server
After=network.target

[Service]
Type=forking
User=new-msi
WorkingDirectory=/home/new-msi/workspace/fsq
ExecStart=/home/new-msi/workspace/fsq/scripts/start-main-server.sh
ExecStop=/home/new-msi/workspace/fsq/scripts/stop-server.sh
Restart=always
RestartSec=10
StandardOutput=append:/home/new-msi/workspace/fsq/logs/systemd-main.log
StandardError=append:/home/new-msi/workspace/fsq/logs/systemd-main.log

[Install]
WantedBy=multi-user.target